# device/systemd/farmadisplay-network.service
[Unit]
Description=TurnoTec Network Healing Daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/turnotec/scripts/network_healing_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
