# device/systemd/farmadisplay-bt-config.service
[Unit]
Description=FarmaDisplay Bluetooth WiFi Configuration Server
After=bluetooth.target
Requires=bluetooth.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/farmadisplay/scripts/bt_wifi_config_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
