# deployment/systemd/farmadisplay-celery.service
[Unit]
Description=FarmaDisplay Celery Worker
After=network.target redis.service
Requires=redis.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/opt/farmadisplay/backend
Environment="PATH=/opt/farmadisplay/backend/venv/bin"
ExecStart=/opt/farmadisplay/backend/venv/bin/celery -A app.celery worker --loglevel=info --detach
ExecStop=/opt/farmadisplay/backend/venv/bin/celery -A app.celery control shutdown
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=farmadisplay-celery

[Install]
WantedBy=multi-user.target
