# device/setup/systemd/turnotec-hotspot.service
[Unit]
Description=TurnoTec Hotspot Manager
After=network.target systemd-rfkill.service
Wants=network.target
Before=turnotec-monitor.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/turnotec/scripts/hotspot_manager.sh auto
ExecStop=/opt/turnotec/scripts/hotspot_manager.sh stop
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=30
StartLimitBurst=3
StartLimitIntervalSec=300

[Install]
WantedBy=multi-user.target
