# device/setup/systemd/turnotec-web.service
[Unit]
Description=TurnoTec Flask Web Server
After=turnotec-hotspot.service
Requires=turnotec-hotspot.service
PartOf=turnotec-hotspot.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/turnotec/web
Environment="FLASK_APP=app.py"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/turnotec/scripts/start_web_server.sh
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

[Install]
WantedBy=multi-user.target
