# device/setup/systemd/turnotec-web.service
[Unit]
Description=TurnoTec Flask Web Server
After=network.target turnotec-hotspot.service
Wants=network.target
PartOf=turnotec-hotspot.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/turnotec/web
Environment="FLASK_APP=app.py"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /opt/turnotec/web/app.py
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
