# device/setup/systemd/turnotec-web.service
[Unit]
Description=TurnoTec Flask Web Server
After=turnotec-hotspot.service
Requires=turnotec-hotspot.service
PartOf=turnotec-hotspot.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/turnotec/web
Environment="FLASK_APP=app.py"
Environment="PYTHONUNBUFFERED=1"
# Stop lighttpd (FullPageOS default web server) to free port 80
ExecStartPre=-/bin/systemctl stop lighttpd
ExecStartPre=-/bin/systemctl disable lighttpd
# Wait 10 seconds after hotspot is ready before starting Flask
ExecStartPre=/bin/sleep 10
ExecStart=/usr/bin/python3 /opt/turnotec/web/app.py
StandardOutput=append:/var/log/turnotec-web.log
StandardError=append:/var/log/turnotec-web.log
Restart=on-failure
RestartSec=15
StartLimitBurst=5
StartLimitIntervalSec=300
TimeoutStartSec=90

[Install]
WantedBy=multi-user.target
